# api/src/config/supervisor.conf
[supervisord]
user = root

[program:app]
command = uwsgi --ini /usr/local/api_service/config/app.ini

[program:nginx]
command = service nginx restart
